# eth-phase-meter.service
[Unit]
Description=ETH Phase Meter - 每4小时自动分析
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/eth-phase-meter.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/eth_phase_meter.py --daemon --interval=4
WorkingDirectory=/root/.openclaw/workspace
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
